# cron/coconutclaw.crontab.example
# Keep agent alive (every minute watchdog)
* * * * * cd "$HOME/code/CoconutClaw" && pgrep -f "coconutclaw run" >/dev/null || nohup coconutclaw run >> LOGS/cron-agent.log 2>&1 &

# Daily heartbeat at 09:00 local time
0 9 * * * cd "$HOME/code/CoconutClaw" && ./scripts/heartbeat.sh >> LOGS/cron-heartbeat.log 2>&1

# Nightly reflection at 22:30 local time
30 22 * * * cd "$HOME/code/CoconutClaw" && ./scripts/nightly_reflection.sh >> LOGS/cron-reflection.log 2>&1
