# systemd/coconutclaw.service
[Unit]
Description=CoconutClaw Telegram Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/code/CoconutClaw
Environment=INSTANCE_DIR=%h/code/CoconutClaw
ExecStart=/usr/bin/env coconutclaw run
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
