# systemd/coconutclaw-heartbeat.service
[Unit]
Description=CoconutClaw Daily Heartbeat

[Service]
Type=oneshot
WorkingDirectory=%h/code/CoconutClaw
Environment=INSTANCE_DIR=%h/code/CoconutClaw
ExecStart=/usr/bin/env coconutclaw heartbeat
