# systemd/coconutclaw-nightly-reflection.service
[Unit]
Description=CoconutClaw Nightly Reflection

[Service]
Type=oneshot
WorkingDirectory=%h/code/CoconutClaw
EnvironmentFile=%h/code/CoconutClaw/.env
ExecStart=/usr/bin/env coconutclaw nightly-reflection
