# systemd/coconutclaw-nightly-reflection.timer
[Unit]
Description=Run CoconutClaw nightly reflection daily

[Timer]
OnCalendar=*-*-* 22:30:00
Persistent=true
Unit=coconutclaw-nightly-reflection.service

[Install]
WantedBy=timers.target
